# systemd/misfans.service
[Unit]
Description=misfans daemon - Raspberry Pi fan controller
After=network.target

[Service]
Type=simple
User=misfans
Group=gpio
EnvironmentFile=/etc/misfans.env
WorkingDirectory=/opt/misfans
ExecStart=/opt/misfans/venv/bin/python -m misfans.daemon
Restart=on-failure
RestartSec=5
# Sandboxing
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=yes
PrivateDevices=yes

[Install]
WantedBy=multi-user.target
